# bebop2/uavpal/conf/70-huawei-e3372.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="12d1", RUN+="/data/ftp/uavpal/bin/uavpal_bebop2.sh"
