# bebop2/uavpal/conf/70-huawei-e3372.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="12d1", RUN+="/usr/bin/flock -n /tmp/lock/uavpal_bebop2 /data/ftp/uavpal/bin/uavpal_bebop2.sh"